# book_management.service
[Service]
ExecStart=/root/bookManagement/run_socket
Restart=always
RestartSec=2s500ms
[Install]
WantedBy=default.target
